# etc/sync.conf
[program:mega-sync]
command=bash -c '/dockerstartup/sync.sh'

directory=/home/www-data
user=root
autorestart=false
startretries=3
