# etc/sync.conf
[program:mega-sync]
command=/dockerstartup/sync.sh

user=www-data
directory=/home/www-data
autorestart=false
startretries=3
