# etc/sync.conf
[program:mega-sync]
command=bash -c '/dockerstartup/sync.sh'

autorestart=false
user=www-data
directory=/home/www-data
