# etc/sync.conf
[program:mega-sync]
command=bash -c '/dockerstartup/sync.sh'

autorestart=false

environment=HOME="/home/www-data",USER="www-data"
