# etc/sync.conf
[program:mega-sync]
command=bash -c '/dockerstartup/sync.sh'

autorestart=true
startretries=3
user=www-data
directory=/home/www-data
