# sv/sync.conf
[program:mega-sync]
command=bash -c '/dockerstartup/startup.sh'

autorestart=true
directory=/root